# dockerfiles/gunicorn.conf
[program:gunicorn]
directory=/app
command=/usr/local/bin/gunicorn --workers 3 --bind unix:/app/app.sock backend.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[group:guni]
programs:gunicorn
